# slapme.conf
# slapme.website websocket stuff

description	"slapme.website websocket server"
author "Paul Harper <benekastah@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

chdir /home/ubuntu/slapme.website
script
    mkdir -p logs
    LOG_FILE="logs/$(date -I).log"
    venv/bin/python -m slapme > "$LOG_FILE" 2> "$LOG_FILE"
end script
